# scrobbler/com.needledrop.scrobbler.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.needledrop.scrobbler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/NeedleDrop/scrobbler/.venv/bin/needledrop</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/NeedleDrop/scrobbler</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/.config/needledrop/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>~/.config/needledrop/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
